# assets/hypr-overlay.service
[Unit]
Description=Discord voice overlay for Hyprland
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/hypr-overlay-wl
Restart=on-failure
RestartSec=3s
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=RUST_LOG=hypr_overlay_wl=debug

[Install]
WantedBy=graphical-session.target
